# etc/systemd/user/kranewl-session.target
[Unit]
Description=kranewl compositor session
BindsTo=graphical-session.target
Wants=graphical-session-pre.target
After=graphical-session-pre.target
